# server/logMasterServer.service
[Unit]
Description=logMaster Server - Journal Log Viewer Web App
Documentation=https://github.com/your-org/logMaster
After=network.target logMasterDaemon.service
Wants=logMasterDaemon.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/logMaster/server
EnvironmentFile=/root/logMaster/server/.env
ExecStart=/root/logMaster/server/.venv/bin/gunicorn -c gunicorn.conf.py logMasterServer:app
Restart=on-failure
RestartSec=5
SyslogIdentifier=logMaster-server

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/logMaster/server

[Install]
WantedBy=multi-user.target
