# daemon/logMasterDaemon.service
[Unit]
Description=logMaster Daemon - journald Log API
Documentation=https://github.com/your-org/logMaster
After=network.target systemd-journald.service
Wants=systemd-journald.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/logMaster/daemon
EnvironmentFile=/root/logMaster/daemon/.env
ExecStart=/root/logMaster/daemon/.venv/bin/gunicorn -c gunicorn.conf.py daemon:app
Restart=on-failure
RestartSec=5

# The daemon needs to read the journal. Add logMaster user to the
# systemd-journal group (see setup.sh) instead of running as root.
SupplementaryGroups=systemd-journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/logMaster/daemon

[Install]
WantedBy=multi-user.target
